# deploy/alex-memory.service
[Unit]
Description=alex-memory MCP server
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=alex
WorkingDirectory=/home/alex/dev/alex-memory
ExecStart=/usr/bin/gleam run
Restart=on-failure
RestartSec=5
Environment=HOME=/home/alex

[Install]
WantedBy=multi-user.target
